# debmetrics/crontab
*/5	*	*	*	*	./pull_runner.py
*/5	*	*	*	*	./graph_runner.py
