# debmetrics/crontab
*/5	*	*	*	*	/var/www/debmetrics/debmetrics/pull_runner.py
